# ansible/roles/monmon/templates/alert.timer.j2
[Unit]
Description=Timer for alert.py service

[Timer]
OnBootSec=70

[Install]
WantedBy=timers.target
