# ansible/roles/monmon/templates/check.timer.j2
[Unit]
Description=Timer for check.py service

[Timer]
OnBootSec=60
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
